# com.chatterbox.notification-reader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chatterbox.notification-reader</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/davidirvine/Desktop/Devel/projects/chatterbox/notification_reader.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/davidirvine/Desktop/Devel/projects/chatterbox</string>
    
    <key>StandardOutPath</key>
    <string>/Users/davidirvine/Desktop/Devel/projects/chatterbox/logs/notification-reader.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/davidirvine/Desktop/Devel/projects/chatterbox/logs/notification-reader.error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CHATTERBOX_URL</key>
        <string>http://127.0.0.1:8000</string>
    </dict>
</dict>
</plist>